# born2beroot/crontab
* * * * * /home/minjkim2/monitoring.sh | wall
